# db-sync.timer
[Unit]
Description=Run KeePassXC Database sync periodically

[Timer]
OnCalendar=Sun
Persistent=true

[Install]
WantedBy=timers.target
